# pi-announce/announce-ip.service
[Unit]
Description=Announce IP address via Pushover
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/announce-ip.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
